# cf-dns-update.service
[Unit]
Description=Cloudflare DNS Updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cf-dns-update
